# templates/default/systemv.erb
#!/bin/bash
# chkconfig: 345 20 80
# description: Play start/shutdown script
# processname: <%= @servicename %>
#
# Installation:
# copy file to /etc/init.d/<%= @servicename %>
# chmod +x /etc/init.d/<%= @servicename %>
# chkconfig --add /etc/init.d/<%= @servicename %>
# chkconfig <%= @servicename %> on
#
# Usage: (as root)
# service <%= @servicename %> start
# service <%= @servicename %> stop
# service <%= @servicename %> status
#
# Remember, you need python 2.6 to run the play command, it doesn't come standard with RedHat/Centos 5.5
# Also, you may want to temporarily remove the >/dev/null for debugging purposes

# Path to play install directory
PLAY_DIST_HOME="<%= @source %>"
PLAY="${PLAY_DIST_HOME}/bin/<%= @executable %>"

# User running the Play process
USER="<%= @user %>"

PID_PATH="<%= @pid_path %>"
PID_FILE="${PID_PATH}/<%= @servicename %>.pid"

if [[ ! -e ${PID_PATH} ]]; then
  mkdir -p ${PID_PATH}
  chown ${USER}:${USER} ${PID_PATH}
fi

CONFIG_FILE="<%= @conf_path %>"
APP_ARGS="<%= @args %>"

# source function library
. /etc/init.d/functions
RETVAL=0

start() {
  echo -n "Starting <%= @servicename %> service: "
  su -s /bin/sh ${USER} -c "( ${PLAY} -Dpidfile.path=${PID_FILE} -Dconfig.file=${CONFIG_FILE} ${APP_ARGS} & )"
  RETVAL=$?

  if [ ${RETVAL} -eq 0 ]; then
    echo_success
  else
    echo_failure
  fi
  echo
}
stop() {
  echo -n "Shutting down Play service: "
  kill $(cat ${PID_FILE})
  RETVAL=$?

  if [ ${RETVAL} -eq 0 ]; then
    echo_success
  else
    echo_failure
  fi
  echo
}
status() {
  ps $(cat ${PID_FILE})
  RETVAL=$?
}
clean() {
  rm -fr ${PID_FILE}
}
case "$1" in
  start)
    clean
    start
    ;;
  stop)
    stop
    ;;
  restart|reload)
    stop
    sleep 10
    start
    ;;
  status)
    status
    ;;
  clean)
    clean
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
esac
exit 0
